# config/files/usr/etc/systemd/user/waybar.service
[Unit]
Description=Waybar
After=graphical-session.service

[Service]
ExecStart=/usr/bin/waybar

[Install]
WantedBy=hyprland-session.target
